# related/auto-build/files/stage3/cde-ldap.socket
# /etc/systemd/system/cde-ldap.socket
[Unit]
Description=Socket for CdE LDAP Server

[Socket]
ListenStream=636
