# related/auto-build/files/stage3/cde-ldap-dummy.service
# /etc/systemd/system/cde-ldap.service
[Unit]
Description=CdE LDAP Dummy for auto-build

[Service]
Type=oneshot
ExecStart=/bin/true
