# related/auto-build/files/stage3/cdedb-app.socket
[Unit]
Description=CdEDB Application server socket

[Socket]
ListenStream=8998
SocketUser=www-data

[Install]
WantedBy=sockets.target
