# related/auto-build/files/stage3/cdedb-app.service
[Unit]
Description=CdEDB Application server
Requires=cdedb-app.socket
After=network.target

[Service]
Type=notify
User=www-cde
Group=www-data
RuntimeDirectory=gunicorn
WorkingDirectory=/cdedb2
Environment="SCRIPT_NAME=/db"
ExecStart=/usr/bin/gunicorn --forwarded-allow-ips="*" -w 4 --reload wsgi.cdedb-app:application
ExecReload=/bin/kill -s HUP $MAINPID
KillMode=mixed
TimeoutStopSec=5
PrivateTmp=true
